# crontab
*/3 * * * * python3 /build/ufos_autotest.py conf -c /build/stand_configs/conf_eb-exp-demo-ufos-jetty.json