# crontab
*/10 * * * * python3 /build/ufos_autotest.py conf -c /build/stand_configs/conf_host1.json
